# scripts/systemd/daqontrol/chips-daqontrol.service.in
[Unit]
Description=DAQ Control

[Service]
ExecStart=%BIN_PATH%/run_with_env.sh daqontrol
WorkingDirectory=%RUN_PATH%
Restart=always
User=%RUN_USER%
Group=%RUN_USER%

[Install]
WantedBy=network.target

